# web-task-automator.service
[Unit]
Description=Web Task Automator Service
After=network.target

[Service]
Type=simple
User=automator
Group=automator
WorkingDirectory=/opt/web-task-automator
Environment=PATH=/opt/web-task-automator/venv/bin
ExecStart=/opt/web-task-automator/venv/bin/python start_api.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/web-task-automator

# Resource limits
LimitNOFILE=65536
MemoryMax=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target

